# system_files/usr/lib/systemd/system/ublue-fix-hostname.service
[Unit]
Description=Fix hostname
After=local-fs.target
ConditionPathExists=!/var/lib/ublue-fix-hostname.done

[Service]
Type=oneshot
ExecStart=/usr/bin/hostnamectl hostname "goose"
ExecStartPost=/usr/bin/touch /var/lib/ublue-fix-hostname.done
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
